# telegram bot/butterfly-history-bot.service.example
[Unit]
Description=Butterfly History Telegram Bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/butterfly-history/telegram bot
ExecStart=/usr/bin/node bot.mjs
Restart=always
RestartSec=5
User=root
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
